# ssh-tunnel-agent.service
[Unit]
Description=SSH Tunnel Agent - Maintains SSH tunnels in a tmux session

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/ssh-tunnel-agent.tmux start
ExecStop=/usr/local/bin/ssh-tunnel-agent.tmux stop
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin

[Install]
WantedBy=default.target
