# ssh-tunnel-agent.service
[Unit]
Description=SSH Tunnel Agent - Maintains SSH tunnels in a tmux session

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=ssh-tunnel-agent.tmux start
ExecStop=ssh-tunnel-agent.tmux stop

[Install]
WantedBy=default.target
