# ssh-tunnel-agent-watcher.service
[Unit]
Description=Restart SSH Tunnel Agent on config change

[Service]
Type=oneshot
ExecStart=systemctl --user restart ssh-tunnel-agent.service
